# config/etc/systemd/system/nftables.service.d/homeserver.conf
[Service]
ExecStart=
ExecStart=/sbin/nft -f /usr/local/share/nftables/main.conf
ExecReload=
ExecReload=/sbin/nft 'flush ruleset; include "/usr/local/share/nftables/main.conf";'
