# config/etc/systemd/system/container-image-builder.timer
[Unit]
Description=Nightly container image builder jobs

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
