# config/etc/systemd/system/gitmirror.timer
[Unit]
Description=Starts gitmirror

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
