# config/etc/systemd/system/googledrive-backup.timer
[Unit]
Description=Starts googledrive backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
