# config/etc/systemd/system/rclone-backup.timer
[Unit]
Description=Starts rclone backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
