# config/etc/systemd/system/restic-autobackup.timer
[Unit]
Description=Starts restic autobackup

[Timer]
OnCalendar=weekly
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
