# config/etc/systemd/system/syncthing-rsync.path
[Unit]
Description=Triggers rsync on syncthing data changes

[Path]
PathChanged=/tmp/syncthing-files-changed

[Install]
WantedBy=multi-user.target
