# files/containers/immich-redis.container
[Unit]
Description=immich-redis
After=network.target

[Container]
ContainerName=immich-redis
Image=registry.hub.docker.com/library/redis:6.2-alpine@sha256:51d6c56749a4243096327e3fb964a48ed92254357108449cb6e23999c37773c5
Network=podman:mac=00:16:3e:42:42:27

[Service]
Restart=always
RestartMode=direct
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
