# config/etc/systemd/system/NetworkManager.service.d/homeserver.conf
[Unit]
After=r8127-modules-load.service
Wants=r8127-modules-load.service
